# obmc-flash-host-cpld@.service.in
[Unit]
Description=Flash Host CPLD image %I to Host

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=echo Please add custom command for flashing image /tmp/image/%i
